# 04Cuarto/Desarrollo_de_aplicaciones_de_internet_DAI/Practica3_semanas_8-9/src/production/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/home/vagrant/gunicorn/socket

[Install]
WantedBy=sockets.target 
